# subsonic.conf
[program:start]
user = nobody
command = /var/subsonic/start.sh